# .config/systemd/user/hyprland-session.target
[Unit]
Description=User systemd services for Hyprland graphical session
#Wants=
BindsTo=graphical-session.target
