# .config/systemd/user/polkit-gnome-authentication-agent.service
[Unit]
Description=Polkit GNOME authentication agent (for sudo GUI)
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
Restart=on-failure

[Install]
WantedBy=graphical-session.target
